# config/claude-governor.service
[Unit]
Description=Claude Governor - automated capacity governor for Claude Code
Documentation=https://github.com/anthropics/claude-governor
After=network-online.target
Wants=network-online.target claude-token-collector.service

[Service]
Type=simple
ExecStart=%h/.local/bin/cgov _daemon
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=default.target
